# files/etc/cron.d/dropboxcheck
# /etc/cron.d/dropboxcheck
#  This deletes downloaded dropbox push-updates

07,37 *     * * *     root   /bin/bash /usr/share/dropbox/dropboxcheck

